# config/khal-remove.service
[Unit]
Description=khal-remove
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/workspace/app
ExecStart=/workspace/app/run

[Install]
WantedBy=multi-user.target
