# config/khal-remove.service
[Unit]
Description=khal-remove

[Service]
WorkingDirectory=/workspace/app
ExecStart=/workspace/app/run

[Install]
WantedBy=multi-user.target
